# config/crontab
0 */3 * * * /srv/update.sh   # every 3 hours
0 0   * * * /srv/droplogs.sh # every day at midnight